# etc/cron.hourly/webpage-cronjob.sh
#!/bin/bash
curl http://www.dal.com
curl http://www.dal.com/academics/programs.html
curl http://www.dal.com/admissions/contactus.html
curl http://www.dal.com/admissions/campus-tours.html
curl http://www.dal.com/academics/requestinfo/domestic-viewbook.html
curl http://www.dal.com/about-dal/leadership-and-vision.html
curl http://www.dal.com/academics/important_dates.html
curl http://www.dal.com/news.html
curl http://www.dal.com/news/events.html
curl http://www.dal.com/cultureofrespect.html

